# recipes-middleware/zenoh/files/zenohd.service
[Unit]
Description=Zenoh Message Router

[Service]
Restart=on-failure
ExecStart=/usr/bin/zenohd

[Install]
WantedBy=multi-user.target
