# recipes-middleware/maivin/files/maivin.target
[Unit]
Description=Maivin Middleware Services
After=multi-user.target
AllowIsolate=yes
Requires=multi-user.target