# recipes-bsp/update-overlays/files/update-overlays.service
[Unit]
Description=update overlays in /boot from /etc/overlays.txt

[Service]
Type=oneshot
ExecStart=/usr/bin/update-overlays
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
